# scripts/aipager.service.example
# Example systemd unit for aipager. Install as a user unit:
#   mkdir -p ~/.config/systemd/user
#   cp aipager.service.example ~/.config/systemd/user/aipager.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now aipager

[Unit]
Description=AIPager Telegram Bot Daemon
After=network-online.target

[Service]
Type=simple
ExecStartPre=-/bin/rm -f /tmp/aipager.sock
ExecStart=%h/.local/bin/aipager start
EnvironmentFile=%h/.config/aipager/config.env
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
